# templates/default/verdaccio.conf.erb
# Ubuntu upstart file at /etc/init/creation-service.conf
# Automatically Generated by Chef

description     "Verdaccio NPM cache & private repo server"

start on (starting network-interface
          or starting network-manager
          or starting networking)

stop on runlevel [!023456]

setuid <%= node['verdaccio']['user'] %>
setgid <%= node['verdaccio']['user'] %>

respawn

script
    exec verdaccio -c <%= File.join(node['verdaccio']['confdir'], 'config.yaml') %>
end script
