# templates/default/verdaccio.service.erb
# Automatically Generated by Chef

[Unit]
Description=Verdaccio NPM cache & private repo server

[Service]
ExecStart=/usr/bin/verdaccio -c <%= File.join(node['verdaccio']['confdir'], 'config.yaml') %>
Type=simple
User=verdaccio
Group=verdaccio

[Install]
WantedBy=multi-user.target
